# automation/com.raj.projectmanager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.raj.projectmanager</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/raj/Desktop/CampusEats/Backend/venv/bin/python3</string>
        <string>/Users/raj/Desktop/CampusEats/automation/project_manager.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardErrorPath</key>
    <string>/Users/raj/Desktop/CampusEats/automation/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/raj/Desktop/CampusEats/automation/output.log</string>
</dict>
</plist>